# normal_hp/normal_hp.conf
[program:normal_hp]
command=python /opt/normal_hp/normal_hp.py 
directory=/opt/normal_hp/
stdout_logfile=/dev/null
stderr_logfile=/dev/null
autostart=true
autorestart=true
startsecs=5
